# main/server/deploy/flutter-reader-api.service
[Unit]
Description=Flutter Reader FastAPI (Uvicorn)
After=network.target

[Service]
Type=simple
# Change to a dedicated user (e.g., app) if created
User=www-data
WorkingDirectory=/opt/flutter_reader/server
EnvironmentFile=/opt/flutter_reader/server/.env
ExecStart=/opt/flutter_reader/venv/bin/uvicorn server.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
